# airootfs/opt/zenkai-installer/zenkai-installer.service
[Unit]
Description=FragOS System Installer.

[Service]
Type=simple
ExecStart=/bin/bash zenkai-installer

[Install]
WantedBy=multi-user.target
